# deploy/systemd/prediction-fund-edge-discovery.timer
[Unit]
Description=Run prediction-fund-edge-discovery 3x daily (9 AM, noon, 3 PM ET)

[Timer]
OnCalendar=*-*-* 13:00:00
OnCalendar=*-*-* 16:00:00
OnCalendar=*-*-* 19:00:00
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
